# config/etc/supervisor/conf.d/mqtt_recorder.conf
[program:mqtt_recorder]
directory=/home/srv/mqtt_recorder/
command=/home/srv/mqtt_recorder/.venv/bin/python run.py
autorestart=true
redirect_stderr=true
